# content/service/metube/run
#!/bin/sh

# Configure & run metube
mkdir -p /mnt/data/videos 2>/dev/null
OUTPUT_TEMPLATE="${YTDL_OUTPUT_TEMPLATE}"
cd /app
exec 2>&1
exec python3 /app/app/main.py